# ansible/roles/etcd/templates/etcd.service.j2
# Generated by Ansible - Do not edit manually
[Unit]
Description=etcd key-value store
Documentation=https://github.com/etcd-io/etcd
After=network.target

[Service]
Type=notify
EnvironmentFile={{ etcd_config_dir }}/etcd.env
ExecStart={{ etcd_bin_dir }}/etcd
Restart=always
RestartSec=10s
LimitNOFILE=65536
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
